# ledmatrix.service
[Unit]
Description=User ledmatrix script
After=multi-user.target

[Service]
Type=simple
ExecStart=~/startup.sh

[Install]
WantedBy=multi-user.target
